# systemd-templates/redditview-tui.service
[Unit]
Description=RedditView Terminal User Interface (TUI)
Documentation=https://github.com/nicthegarden/redditiew
After=network.target redditview-api.service
Requires=redditview-api.service

[Service]
Type=forking
User=__USERNAME__
WorkingDirectory=__INSTALL_PATH__

# Start tmux server in background with proper environment
ExecStartPre=__TMUX_BIN__ new-session -d -s redditview -x 120 -y 40 -c __INSTALL_PATH__

# Run TUI application in the tmux session with proper TTY
ExecStart=__TMUX_BIN__ send-keys -t redditview:0 "__TUI_BIN__" Enter

# Environment variables
Environment="DISPLAY=:0"
Environment="XAUTHORITY=%h/.Xauthority"

# Auto-restart on failure
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redditview-tui

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

# Stop script to clean up tmux session
ExecStopPost=/bin/bash -c "__TMUX_BIN__ kill-session -t redditview || true"

[Install]
WantedBy=multi-user.target
